# examples/zram/zram-writeback-perl-1.0.0/systemd/zram-writeback-budget.timer
[Unit]
Description=Daily zram writeback budget reset
Documentation=man:zram-writeback(8)

[Timer]
OnCalendar=daily
AccuracySec=5min
Persistent=true
Unit=zram-writeback-budget.service

[Install]
WantedBy=timers.target
