# d-i/debian/hooks/shared/target/etc/systemd/system/wpa_supplicant.service.d/10-preseed-no-p2p.conf
# Managed by debian-preseed-di.
#
# Keep the D-Bus service enabled for NetworkManager and desktop clients, but do
# not start a dedicated P2P device interface. The managed main config disables
# Wi-Fi Direct/P2P while preserving normal WPA/SAE station use.

[Service]
ExecStart=
ExecStart=/usr/sbin/wpa_supplicant -u -s -c /etc/wpa_supplicant/wpa_supplicant.conf -O "DIR=/run/wpa_supplicant GROUP=netdev"
