# d-i/debian/hooks/role/desktop/target/etc/systemd/system/bluetooth.service.d/10-preseed-directory-mode.conf
[Service]
# Debian's bluez package ships /etc/bluetooth as 0755. Match that mode so
# systemd does not warn that ConfigurationDirectoryMode differs at start.
ConfigurationDirectoryMode=0755
