# d-i/debian/hooks/shared/target/etc/systemd/system/fstrim.timer.d/override.conf
[Timer]
# Clear vendor weekly schedule
OnCalendar=

# Run every 4 days after the last successful activation
OnBootSec=20min
OnUnitActiveSec=4d

# Avoid exact synchronized execution and reduce boot/load spikes
AccuracySec=1h
RandomizedDelaySec=30min

# Harmless here; useful if you later switch back to calendar scheduling
Persistent=true
