# d-i/debian/hooks/shared/target/etc/systemd/system/bootprofile-apply.service.tmpl
[Unit]
Description=Apply boot profile tunables
DefaultDependencies=no
After=systemd-remount-fs.service
Before=systemd-sysctl.service sysinit.target
ConditionPathExists=/proc/cmdline

[Service]
Type=oneshot
ExecStart=__INSTALLER_FILE_BOOTPROFILE_APPLY__

[Install]
WantedBy=sysinit.target
