# d-i/debian/hooks/shared/target/etc/systemd/system/firstboot.service
[Unit]
Description=First boot provisioning and validation wrapper
Documentation=man:systemd.service(5)
DefaultDependencies=no

Requires=local-fs.target
After=local-fs.target systemd-journald.socket
Before=sysinit.target

ConditionPathExists=!/var/lib/preseed/firstboot/complete

Conflicts=shutdown.target
Before=shutdown.target

[Service]
Type=simple
ExecStart=/usr/local/sbin/firstboot.sh

StandardOutput=journal+console
StandardError=journal+console

TimeoutStartSec=0
TimeoutStopSec=30s
KillSignal=SIGTERM

User=root
Group=root
UMask=0077

[Install]
WantedBy=sysinit.target
