# d-i/debian/hooks/shared/target/etc/systemd/system/grub-btrfs-refresh.path
[Unit]
Description=Watch Timeshift runtime state and refresh the GRUB snapshot menu
ConditionPathExists=/usr/local/sbin/grub-btrfs-refresh

[Path]
PathExistsGlob=/run/timeshift/*
PathChanged=/run/timeshift
Unit=grub-btrfs-refresh.service

[Install]
WantedBy=multi-user.target
