# d-i/debian/hooks/shared/target/etc/systemd/system/podman-rootless-linger.service.tmpl
[Unit]
Description=Enable lingering for managed rootless Podman service user
Documentation=man:loginctl(1)
After=systemd-user-sessions.service
Wants=systemd-user-sessions.service
ConditionPathExists=!__INSTALLER_PODMAN_LINGER_MARKER__

[Service]
Type=oneshot
ExecStart=/usr/bin/loginctl enable-linger __INSTALLER_PODMAN_SERVICE_USER__
ExecStart=/usr/bin/systemctl start user@__INSTALLER_PODMAN_SERVICE_UID__.service
ExecStart=/bin/sh -eu -c 'i=0; while [ ! -S /run/user/__INSTALLER_PODMAN_SERVICE_UID__/bus ]; do i=$((i + 1)); [ "$i" -lt 31 ] || exit 1; sleep 1; done'
ExecStart=/usr/sbin/runuser -u __INSTALLER_PODMAN_SERVICE_USER__ -- /usr/bin/env HOME=__INSTALLER_PODMAN_SERVICE_HOME__ XDG_RUNTIME_DIR=/run/user/__INSTALLER_PODMAN_SERVICE_UID__ DBUS_SESSION_BUS_ADDRESS=unix:path=/run/user/__INSTALLER_PODMAN_SERVICE_UID__/bus PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin /usr/bin/systemctl --user daemon-reload
ExecStart=/usr/sbin/runuser -u __INSTALLER_PODMAN_SERVICE_USER__ -- /usr/bin/env HOME=__INSTALLER_PODMAN_SERVICE_HOME__ XDG_RUNTIME_DIR=/run/user/__INSTALLER_PODMAN_SERVICE_UID__ DBUS_SESSION_BUS_ADDRESS=unix:path=/run/user/__INSTALLER_PODMAN_SERVICE_UID__/bus PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin /usr/bin/systemctl --user start buildah-env.service
ExecStart=/usr/sbin/runuser -u __INSTALLER_PODMAN_SERVICE_USER__ -- /usr/bin/env HOME=__INSTALLER_PODMAN_SERVICE_HOME__ XDG_RUNTIME_DIR=/run/user/__INSTALLER_PODMAN_SERVICE_UID__ DBUS_SESSION_BUS_ADDRESS=unix:path=/run/user/__INSTALLER_PODMAN_SERVICE_UID__/bus PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin /usr/bin/systemctl --user start __INSTALLER_PODMAN_API_START_UNITS__
ExecStart=/usr/bin/install -d -m 0700 /var/lib/preseed/firstboot
ExecStart=/usr/bin/touch __INSTALLER_PODMAN_LINGER_MARKER__
ExecStart=/usr/bin/chmod 0600 __INSTALLER_PODMAN_LINGER_MARKER__

[Install]
WantedBy=multi-user.target
