# d-i/debian/hooks/shared/target/etc/systemd/system/swap-fallback.service.tmpl
[Unit]
Description=Configure encrypted fallback swap partition
After=local-fs.target systemd-modules-load.service
Wants=local-fs.target systemd-modules-load.service
Before=multi-user.target
ConditionPathExists=__INSTALLER_FILE_SWAP_FALLBACK_CONFIG__

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=__INSTALLER_FILE_SWAP_FALLBACK_HELPER__ start
ExecStop=__INSTALLER_FILE_SWAP_FALLBACK_HELPER__ stop
TimeoutStartSec=2min
TimeoutStopSec=1min

[Install]
WantedBy=multi-user.target
