# d-i/debian/hooks/shared/target/etc/systemd/system/timeshift-daily.service
[Unit]
Description=Create the managed four-times-daily Timeshift snapshot
ConditionPathExists=/usr/bin/timeshift
ConditionPathExists=/etc/timeshift/timeshift.json
After=local-fs.target

[Service]
Type=oneshot
Nice=19
IOSchedulingClass=idle
ExecStart=/usr/local/sbin/timeshift-managed-snapshot daily
