# d-i/debian/hooks/shared/target/etc/systemd/system/zram-setup.service.tmpl
[Unit]
Description=Configure zram swap with encrypted writeback backing
After=local-fs.target systemd-modules-load.service
Wants=local-fs.target systemd-modules-load.service
Before=multi-user.target
ConditionPathExists=__INSTALLER_FILE_ZRAM_DEFAULT__
ConditionPathExists=__INSTALLER_FILE_ZRAM_SETUP_HELPER__

[Service]
Type=oneshot
RemainAfterExit=yes
RuntimeDirectory=zram
RuntimeDirectoryMode=0750
ExecStart=__INSTALLER_FILE_ZRAM_SETUP_HELPER__ start
ExecStop=__INSTALLER_FILE_ZRAM_SETUP_HELPER__ stop
TimeoutStartSec=2min
TimeoutStopSec=1min
MemoryAccounting=yes
TasksAccounting=yes
TasksMax=64
OOMPolicy=stop
UMask=0077

[Install]
WantedBy=multi-user.target
