# d-i/debian/hooks/shared/target/etc/systemd/system/zram-writeback.service.tmpl
[Unit]
Description=Run pressure-aware zram recompression/writeback maintenance
After=__INSTALLER_ZRAM_SETUP_UNIT__ sys-kernel-debug.mount
Requires=__INSTALLER_ZRAM_SETUP_UNIT__
Wants=sys-kernel-debug.mount
ConditionPathExists=/sys/block/__INSTALLER_ZRAM_SWAP_DEVICE_NAME__/writeback
ConditionPathIsReadWrite=/sys/block/__INSTALLER_ZRAM_SWAP_DEVICE_NAME__/writeback
ConditionPathExists=__INSTALLER_ZRAM_SWAP_DEVICE__
ConditionPathExists=__INSTALLER_ZRAM_BACKING_DEVICE__
ConditionPathExists=__INSTALLER_FILE_ZRAM_DEFAULT__
ConditionPathExists=__INSTALLER_FILE_ZRAM_CONFIG__
ConditionPathExists=__INSTALLER_FILE_ZRAM_WRITEBACK_HELPER__
StartLimitIntervalSec=30min
StartLimitBurst=2

[Service]
Type=oneshot
EnvironmentFile=__INSTALLER_FILE_ZRAM_DEFAULT__
RuntimeDirectory=zram
RuntimeDirectoryMode=0750
ExecStart=__INSTALLER_FILE_ZRAM_WRITEBACK_HELPER__ --config ${ZRAM_POLICY_CONFIG} run
TimeoutStartSec=5min
Nice=19
CPUSchedulingPolicy=idle
CPUWeight=10
IOSchedulingClass=idle
IOSchedulingPriority=7
IOAccounting=yes
IOWeight=10
IOWriteBandwidthMax=__INSTALLER_ZRAM_BACKING_DEVICE__ __INSTALLER_ZRAM_MAINTENANCE_IO_WRITE_BANDWIDTH_MAX__
MemoryAccounting=yes
MemoryHigh=__INSTALLER_ZRAM_MAINTENANCE_MEMORY_HIGH__
MemoryMax=__INSTALLER_ZRAM_MAINTENANCE_MEMORY_MAX__
OOMScoreAdjust=500
NoNewPrivileges=yes
PrivateTmp=yes
ProtectHome=yes
ProtectControlGroups=yes
ProtectKernelModules=yes
ProtectSystem=strict
ReadWritePaths=__INSTALLER_ZRAM_RUNTIME_DIR__
RestrictRealtime=yes
RestrictSUIDSGID=yes
LockPersonality=yes
MemoryDenyWriteExecute=yes
SystemCallArchitectures=native
UMask=0077
