# d-i/debian/hooks/shared/target/data/config/podman/templates/rootless/systemd/user/podman.service.d/10-podman-service-managed.conf
[Service]
Environment=LOGGING=--log-level=error
UnsetEnvironment=DOCKER_HOST CONTAINER_HOST
Delegate=true
ExecStart=
ExecStart=/usr/bin/podman $LOGGING system service --time=0
__INSTALLER_PODMAN_SERVICE_SLICE_LINE__
