# d-i/debian/hooks/shared/target/data/config/podman/templates/rootless/systemd/user/podman.socket.d/10-podman-socket-managed.conf
[Socket]
SocketMode=0600
DirectoryMode=0700
