# d-i/debian/hooks/shared/target/data/config/podman/templates/podbin/systemd/users/container.container.tmpl
[Unit]
Description=Podbin managed rootless container __PODBIN_CONTAINER_NAME__
Documentation=man:podman-systemd.unit(5)

[Container]
Image=__PODBIN_IMAGE__
ContainerName=__PODBIN_CONTAINER_NAME__
Pull=missing
PublishPort=__PODBIN_BIND_IP__:__PODBIN_HOST_PORT__:__PODBIN_CONTAINER_SSH_PORT__/tcp
ReadOnly=__PODBIN_READ_ONLY__
NoNewPrivileges=true
DropCapability=all
User=__PODBIN_RUNTIME_USER_UID__:__PODBIN_RUNTIME_USER_GID__
Tmpfs=/tmp:rw,mode=1777
Tmpfs=/run/sshd:rw,mode=0755,uid=__PODBIN_RUNTIME_USER_UID__,gid=__PODBIN_RUNTIME_USER_GID__
Tmpfs=/var/tmp:rw,mode=1777
Tmpfs=__PODBIN_RUNTIME_USER_HOME__:rw,mode=0700,uid=__PODBIN_RUNTIME_USER_UID__,gid=__PODBIN_RUNTIME_USER_GID__
Tmpfs=__PODBIN_RUNTIME_WORKDIR__:rw,mode=0755,uid=__PODBIN_RUNTIME_USER_UID__,gid=__PODBIN_RUNTIME_USER_GID__
Volume=__PODBIN_DATA_DIR__/ssh:__PODBIN_AUTH_KEYS_DIR__:ro
EnvironmentFile=__PODBIN_ENV_FILE__
Label=io.debian-preseed-di.podbin.user=__PODBIN_USER_NAME__
Label=io.debian-preseed-di.podbin.container=__PODBIN_CONTAINER__

[Service]
Restart=always
RestartSec=5
TimeoutStartSec=900

[Install]
WantedBy=default.target
