# d-i/debian/hooks/shared/target/data/config/podman/templates/rootless/systemd/user/podman-api-env.service.tmpl
[Unit]
Description=Managed rootless Podman API environment for __INSTALLER_PODMAN_SERVICE_USER__
After=podman.socket
Wants=podman.socket

[Service]
Type=oneshot
RemainAfterExit=yes
__INSTALLER_PODMAN_API_SERVICE_ENVIRONMENT_LINES__
ExecStart=/usr/bin/systemctl --user set-environment __INSTALLER_PODMAN_API_SET_ENV_ARGS__
ExecStop=/usr/bin/systemctl --user unset-environment __INSTALLER_PODMAN_API_UNSET_ENV_NAMES__

[Install]
WantedBy=default.target
