# d-i/debian/hooks/shared/target/etc/systemd/system/unattended-upgrades.service.d/10-preseed-warning-policy.conf
# Managed by debian-preseed-di.
#
# Debian's unattended-upgrade-shutdown helper can emit a PyGObject deprecation
# warning for GLib.unix_signal_add during shutdown. Keep shutdown behavior
# unchanged and suppress deprecation-only noise from this package-owned helper.

[Service]
Environment=PYTHONWARNINGS=ignore::DeprecationWarning
