# d-i/debian/hooks/shared/target/etc/systemd/system/auditd.service.d/override.conf
[Service]
ExecStartPost=
ExecStartPost=-/usr/local/sbin/augenrules-quiet --load
