# d-i/debian/hooks/services/gitlab-runner/target/etc/systemd/system/aptly-bridge.path
[Unit]
Description=Watch for queued Aptly publish bridge requests
ConditionPathExists=/usr/local/sbin/aptly-bridge-processor

[Path]
PathExistsGlob=/pool/aptly/queue/requests/*.json
Unit=aptly-bridge.service

[Install]
WantedBy=multi-user.target
