# d-i/debian/hooks/hardware/cpu/intel/target/etc/udev/rules.d/85-wifi-regdom.rules
ACTION=="add", SUBSYSTEM=="ieee80211", RUN+="/sbin/iw reg set AU"
