# d-i/debian/hooks/shared/target/etc/systemd/system/apt-daily.service.d/override.conf.tmpl
[Service]
Environment=DEBIAN_FRONTEND=noninteractive
Environment=APT_LISTCHANGES_FRONTEND=none
Nice=10
IOSchedulingClass=best-effort
IOSchedulingPriority=7
TimeoutStartSec=30min
PrivateTmp=true
PrivateDevices=true
ProtectSystem=full
ProtectHome=true
ProtectClock=yes
ProtectControlGroups=yes
ProtectHostname=yes
ProtectKernelLogs=yes
ProtectKernelModules=yes
ProtectKernelTunables=yes
NoNewPrivileges=yes
LockPersonality=yes
RestrictSUIDSGID=yes
RestrictRealtime=yes
SystemCallArchitectures=native
